# ops/wifi-watchdog/wifi-watchdog.service
[Unit]
Description=ShlimazlBot Wi-Fi/network watchdog
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wifi-watchdog.sh
