# ops/wifi-watchdog/wifi-watchdog.timer
[Unit]
Description=Run ShlimazlBot Wi-Fi/network watchdog every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=30s
Unit=wifi-watchdog.service

[Install]
WantedBy=timers.target
